# root/etc/profile.d/toolkit.sh
#!/usr/bin/env bash

export PATH="/opt/toolkit/ass/bin:${PATH}"

export XDG_DATA_DIRS="/opt/toolkit/ass:${XDG_DATA_DIRS:-}"

# fcitx5 >/dev/null 2>&1 &
